# com.local-mcp-kb.watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local-mcp-kb.watcher</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/marcostalman/Projects/local-mcp-kb/.venv/bin/python</string>
        <string>/Users/marcostalman/Projects/local-mcp-kb/src/watcher.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/marcostalman/Projects/local-mcp-kb</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/marcostalman/.local-mcp-kb/watcher.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/marcostalman/.local-mcp-kb/watcher-error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
